# Sigma/KarstoRat/KarstoRat_registry_persistence.yml
title: KarstoRat Registry Persistence
id: 63b1d3ef-af7c-4f6d-9b8d-5a0a3c4e8f1a
status: experimental
description: Detects KarstoRat adding persistence via Windows Registry Run key
author: Cyberthint Threat Intelligence Team
date: 2026-03-03
modified: 2026-03-03
references:
    - https://cyberthint.io/karstorat-rat-case-study
    - https://github.com/cyberthint/karstorat-detection-rules
tags:
    - attack.persistence
    - attack.t1547.001
    - malware.karstorat
logsource:
    product: windows
    category: registry_event
detection:
    selection:
        EventType: SetValue
        TargetObject|contains:
            - 'CurrentVersion\Run'
        Details|contains: 'SecurityService.exe'
    condition: selection
falsepositives:
    - Legitimate security software using similar names
    - System administrators deploying tools with similar naming
level: high